# crash.service
[Unit]
Description=Minecraft Reset

[Service]
ExecStart=/var/minecraft/crash.sh

[Install]
WantedBy=multi-user.target
