# crash.service
[Unit]
Description=Minecraft Crash Protection

[Service]
ExecStart=/var/minecraft/crash.sh

[Install]
WantedBy=multi-user.target
